# deploy/systemd/nanobot-update.service
[Unit]
Description=nanobot self-update task
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/nanobot
Environment=APP_DIR=%h/nanobot
Environment=SERVICE_NAME=nanobot.service
Environment=SYSTEMCTL_CMD=systemctl
Environment=SYSTEMCTL_SCOPE=--user
Environment=UV_SYNC_ARGS=--extra tts
ExecStart=/bin/bash %h/nanobot/deploy/update.sh

# Security hardening (optional, relax if your environment needs it)
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=default.target
