# deploy/systemd/nanobot.service
[Unit]
Description=nanobot gateway service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nanobot
Group=nanobot
WorkingDirectory=/opt/nanobot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env uv run nanobot gateway
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=30

# Journald logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
