# deploy/systemd/nanobot-update.service
[Unit]
Description=nanobot self-update task
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/nanobot
ExecStart=/bin/bash %h/nanobot/deploy/update.sh

# Security hardening (optional, relax if your environment needs it)
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=default.target
