# deploy/systemd/nanobot-upgrade.service
[Unit]
Description=nanobot self-upgrade task
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/nanobot
Environment=APP_DIR=/opt/nanobot
Environment=SERVICE_NAME=nanobot.service
ExecStart=/bin/bash /opt/nanobot/deploy/upgrade.sh

# Security hardening (optional, relax if your environment needs it)
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
